# doc/flowintel-notifications.service
[Unit]
Description=Flowintel Notification Service
After=network.target postgresql.service valkey.service
Requires=postgresql.service valkey.service

[Service]
Type=simple
User=yourusername
Group=yourusername
WorkingDirectory=/opt/flowintel/flowintel
Environment="PATH=/opt/flowintel/flowintel/env/bin:/usr/local/bin:/usr/bin:/bin"
Environment="FLASKENV=production"
Environment="HISTORY_DIR=/opt/flowintel/flowintel/history"
ExecStart=/opt/flowintel/flowintel/env/bin/python startNotif.py
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
